# packages/ffopenmppt/files/etc/uci-defaults/90-ffopenmppt-cron
#!/bin/sh

unique_append()
{
	grep -qF "$1" "$2" || echo "$1" >> "$2"
}

unique_append \
	'0 22 * * * echo "0" > /tmp/battery-gauge-stop.log && tail -n 1 /tmp/ISEMS/ffopenmppt.log | cut -d \; -f 10 > /tmp/battery-gauge-start.log'\
	
	/etc/crontabs/root


unique_append \
	'0 4 * * *  tail -n 1 /tmp/ISEMS/ffopenmppt.log | cut -d \; -f 10 > /tmp/battery-gauge-stop.log'\
	/etc/crotabs/root
	
exit 0
