# packages/ffopenmppt/files/etc/init.d/ffopenmppt
#!/bin/sh /etc/rc.common

START=95
STOP=95

# Read serial port setting from /etc/config/ffopenmppt

MPPT_TTY="$(uci get ffopenmppt.@ffopenmppt[0].serial_port)"

boot() {

# Multiplex serial port for reading.

stty -F /dev/$MPPT_TTY 9600 -ixon -ixoff -ixany -echo -echoe -echok -echonl -tostop -echoprt -echoctl -echoke -crtscts

mknod /tmp/SERIAL_0 p
mknod /tmp/SERIAL_1 p 

cat /dev/$MPPT_TTY | tee  /tmp/SERIAL_0 /tmp/SERIAL_1 1&> /dev/null & 

pid=$!

echo $pid > /tmp/serport-mpx-pid

/usr/bin/collect-ISEMS-data.sh & 

pid=$!

echo $pid > /tmp/collect-isems-data-pid

/usr/bin/powersave.sh &

pid=$!

echo $pid > /tmp/powersave-pid

}

stop() {

kill `cat /tmp/serport-mpx-pid`
kill `cat /tmp/collect-isems-data-pid`
kill `cat /tmp/powersave-pid`
rm /tmp/SERIAL_0
rm /tmp/SERIAL_1
rm /tmp/serport-mpx-pid
rm /tmp/collect-isems-data-pid
}
